# 5/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/bin/bash -c "exec /usr/sbin/nginx"

[program:spawn-fcgi]
command=/bin/bash -c "exec  spawn-fcgi -p 8080 /home/server"